# dbus-1.13.6/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/home/jianglei/svg/build/usr/local/var/run/dbus/system_bus_socket
